# examples/elastic/windows_registry_persistence_run_keys.json
{
  "rule_id": "5fda1640-0c25-4c30-9ddb-cf6b92c5af9e",
  "name": "Windows Registry Run Key Persistence",
  "description": "Detects modifications to Windows Registry Run keys which could be used for persistence",
  "author": ["RuleHub Examples"],
  "references": [
    "https://attack.mitre.org/techniques/T1547/001/",
    "https://docs.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys"
  ],
  "severity": "medium",
  "risk_score": 47,
  "rule_type": "query",
  "type": "eql",
  "created_at": "2025-05-16T00:00:00.000Z",
  "updated_at": "2025-05-16T00:00:00.000Z",
  "query": "registry where registry.path : (\"HKEY_USERS\\\\*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\*\", \"HKEY_USERS\\\\*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\*\", \"HKEY_LOCAL_MACHINE\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\*\", \"HKEY_LOCAL_MACHINE\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\*\") and process.name != \"msiexec.exe\"",
  "language": "eql",
  "tags": [
    "Windows", "Persistence", "Registry"
  ],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1547.001",
          "name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder",
          "reference": "https://attack.mitre.org/techniques/T1547/001/"
        }
      ]
    }
  ],
  "false_positives": [
    "Software installation processes",
    "Legitimate system administration tasks",
    "Startup programs configured by administrators"
  ],
  "os_types": ["windows"],
  "index": ["winlogbeat-*", "logs-endpoint.events.*"],
  "related_integrations": [
    {
      "package": "endpoint",
      "version": ">=7.9.0"
    }
  ],
  "required_fields": [
    {
      "ecs": true,
      "name": "process.name",
      "type": "keyword"
    },
    {
      "ecs": true,
      "name": "registry.path",
      "type": "keyword"
    }
  ],
  "setup": "This rule requires Elastic Endpoint Security integration or Winlogbeat with registry monitoring configured."
}