# examples/splunk/linux_suspicious_commands.yml
name: Linux Suspicious Commands Execution
id: 7123a642-e28d-4b6a-9ca9-d841f3f75eb5
version: 1
date: '2025-05-16'
modified: '2025-05-16'
author: RuleHub Examples
description: Detects execution of commands commonly used by attackers for reconnaissance and lateral movement on Linux systems
type: detection
search: >-
  index=linux sourcetype=linux_secure OR sourcetype=linux_audit OR sourcetype=syslog
  (
    "wget http" OR "curl http" OR "chmod +x" OR "/dev/shm/" OR
    "nc -l" OR "netcat -l" OR "python -m SimpleHTTPServer" OR "python -m http.server" OR
    "echo SSH-" OR "masscan" OR "nmap -p" OR "find / -perm" OR
    "base64 -d" OR "ping -c 1" OR "/tmp/..." OR "crontab -e" OR
    "socat" OR "mknod /tmp" OR "mkfifo /tmp" OR "msfvenom" OR
    "id;uname -a" OR "whoami;cat /etc/passwd" OR "chattr -i" OR
    "iptables -F" OR "passwd root" OR "adduser" OR "usermod -aG sudo"
  )
  | stats count min(_time) as firstTime max(_time) as lastTime by host src_user command
  | rename host as dest
  | table firstTime lastTime dest src_user command count
data_source:
  - Linux Audit Logs
  - Linux Secure Logs
  - Syslog
tags:
  - Linux
  - Reconnaissance
  - Lateral Movement
  - Persistence
  - attack.discovery
  - attack.lateral_movement
  - attack.persistence
risk_score: 65
risk_severity: high
mitre_attack_id:
  - T1087
  - T1083
  - T1082
  - T1016
  - T1046
  - T1136
  - T1543.003
  - T1021
correlation_rule: false
known_false_positives:
  - System administrators may use these commands for legitimate purposes
  - DevOps automation scripts
  - Security testing and vulnerability scanning
fields:
  - _time
  - host
  - src_user
  - command
status: production
asset_type: Linux
playbook: >-
  1. Verify if the command execution is legitimate by checking with the system administrator
  2. Check the user account that executed the command to determine if it's authorized
  3. Review surrounding events to identify potential compromise indicators
  4. If suspicious, isolate the system and conduct deeper forensic investigation
response_actions:
  - Isolate affected systems if malicious activity is confirmed
  - Reset compromised user accounts
  - Block malicious IP addresses involved in the activity
  - Scan system for persistence mechanisms and malware
testing:
  sample_data: |
    May 16 12:34:56 web-server sshd[1234]: Accepted password for alice from 192.168.1.100 port 41234
    May 16 12:35:01 web-server bash[1235]: alice: COMMAND=/bin/bash -c id;uname -a
  validation_query: |
    index=linux sourcetype=linux_secure (id;uname -a)