# rules/sigma/1c0e41cd-21bb-4433-9acc-4a2cd6367b9b.json
{
  "id": "1c0e41cd-21bb-4433-9acc-4a2cd6367b9b",
  "name": "Suspicious Modification Of Scheduled Tasks",
  "description": "Detects when an attacker tries to modify an already existing scheduled tasks to run from a suspicious location\nAttackers can create a simple looking task in order to avoid detection on creation as it's often the most focused on\nInstead they modify the task after creation to include their malicious payload\n",
  "source": {
    "type": "sigma",
    "id": "1c0e41cd-21bb-4433-9acc-4a2cd6367b9b",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_change.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_change.yml"
  },
  "tags": [
    "attackexecution",
    "attackt1053005"
  ],
  "author": "Nasreddine Bencherchali (Nextron Systems)",
  "references": [
    "Internal Research",
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2022-07-28",
  "modified": "2022-11-18",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_schtasks: {\"Image|endswith\": \"\\\\schtasks.exe\", \"CommandLine|contains|all\": [\" /Change \", \" /TN \"]}\nselection_susp_locations: {\"CommandLine|contains\": [\"\\\\AppData\\\\Local\\\\Temp\", \"\\\\AppData\\\\Roaming\\\\\", \"\\\\Users\\\\Public\\\\\", \"\\\\WINDOWS\\\\Temp\\\\\", \"\\\\Desktop\\\\\", \"\\\\Downloads\\\\\", \"\\\\Temporary Internet\", \"C:\\\\ProgramData\\\\\", \"C:\\\\Perflogs\\\\\", \"%ProgramData%\", \"%appdata%\", \"%comspec%\", \"%localappdata%\"]}\nselection_susp_images: {\"CommandLine|contains\": [\"regsvr32\", \"rundll32\", \"cmd /c \", \"cmd /k \", \"cmd /r \", \"cmd.exe /c \", \"cmd.exe /k \", \"cmd.exe /r \", \"powershell\", \"mshta\", \"wscript\", \"cscript\", \"certutil\", \"bitsadmin\", \"bash.exe\", \"bash \", \"scrcons\", \"wmic \", \"wmic.exe\", \"forfiles\", \"scriptrunner\", \"hh.exe\", \"hh \"]}",
    "condition": "all of selection_*",
    "fields": []
  },
  "falsepositives": [
    "Unknown"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}