# rules/sigma/24c8392b-aa3c-46b7-a545-43f71657fe98.json
{
  "id": "24c8392b-aa3c-46b7-a545-43f71657fe98",
  "name": "Suspicious Schtasks Schedule Types",
  "description": "Detects scheduled task creations or modification on a suspicious schedule type",
  "source": {
    "type": "sigma",
    "id": "24c8392b-aa3c-46b7-a545-43f71657fe98",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_schedule_type.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_schedule_type.yml"
  },
  "tags": [
    "attackexecution",
    "attackt1053005"
  ],
  "author": "Nasreddine Bencherchali (Nextron Systems)",
  "references": [
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change",
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create",
    "http://blog.talosintelligence.com/2022/09/lazarus-three-rats.html"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2022-09-09",
  "modified": "2022-09-09",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_img: [{\"Image|endswith\": \"\\\\schtasks.exe\"}, {\"OriginalFileName\": \"schtasks.exe\"}]\nselection_time: {\"CommandLine|contains\": [\" ONLOGON \", \" ONSTART \", \" ONCE \", \" ONIDLE \"]}\nfilter_privs: {\"CommandLine|contains\": [\"NT AUT\", \" SYSTEM\", \"HIGHEST\"]}",
    "condition": "all of selection_* and not 1 of filter_*",
    "fields": []
  },
  "falsepositives": [
    "Legitimate processes that run at logon. Filter according to your environment"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}