# rules/sigma/28208707-fe31-437f-9a7f-4b1108b94d2e.json
{
  "id": "28208707-fe31-437f-9a7f-4b1108b94d2e",
  "name": "Suspicious Startup Folder Persistence",
  "description": "Detects when a file with a suspicious extension is created in the startup folder",
  "source": {
    "type": "sigma",
    "id": "28208707-fe31-437f-9a7f-4b1108b94d2e",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/file_event_win_susp_startup_folder_persistence.yml",
    "file_path": "tmp/repos/sigma/rules/windows/file/file_event/file_event_win_susp_startup_folder_persistence.yml"
  },
  "tags": [
    "attackpersistence",
    "attackt1547001"
  ],
  "author": "Nasreddine Bencherchali (Nextron Systems)",
  "references": [
    "https://github.com/last-byte/PersistenceSniper"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2022-08-10",
  "modified": "2023-01-06",
  "mitre": {
    "tactics": [
      "TA1547"
    ],
    "techniques": [
      "T001"
    ]
  },
  "detection": {
    "query": "selection: {\"TargetFilename|contains\": \"\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\\", \"TargetFilename|endswith\": [\".vbs\", \".vbe\", \".bat\", \".ps1\", \".hta\", \".dll\", \".jar\", \".msi\", \".scr\", \".cmd\"]}",
    "condition": "selection",
    "fields": []
  },
  "falsepositives": [
    "Rare legitimate usage of some of the extensions mentioned in the rule"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "file_event"
  ]
}