# rules/sigma/39019a4e-317f-4ce3-ae63-309a8c6b53c5.json
{
  "id": "39019a4e-317f-4ce3-ae63-309a8c6b53c5",
  "name": "Suspicious Scheduled Task Creation Involving Temp Folder",
  "description": "Detects the creation of scheduled tasks that involves a temporary folder and runs only once",
  "source": {
    "type": "sigma",
    "id": "39019a4e-317f-4ce3-ae63-309a8c6b53c5",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_creation_temp_folder.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_creation_temp_folder.yml"
  },
  "tags": [
    "attackexecution",
    "attackpersistence",
    "attackt1053005"
  ],
  "author": "Florian Roth (Nextron Systems)",
  "references": [
    "https://discuss.elastic.co/t/detection-and-response-for-hafnium-activity/266289/3"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2021-03-11",
  "modified": "2022-10-09",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection: {\"Image|endswith\": \"\\\\schtasks.exe\", \"CommandLine|contains|all\": [\" /create \", \" /sc once \", \"\\\\Temp\\\\\"]}",
    "condition": "selection",
    "fields": []
  },
  "falsepositives": [
    "Administrative activity",
    "Software installation"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}