# rules/sigma/3a734d25-df5c-4b99-8034-af1ddb5883a4.json
{
  "id": "3a734d25-df5c-4b99-8034-af1ddb5883a4",
  "name": "Suspicious Scheduled Task Creation",
  "description": "Detects suspicious scheduled task creation events. Based on attributes such as paths, commands line flags, etc.",
  "source": {
    "type": "sigma",
    "id": "3a734d25-df5c-4b99-8034-af1ddb5883a4",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/win_security_susp_scheduled_task_creation.yml",
    "file_path": "tmp/repos/sigma/rules/windows/builtin/security/win_security_susp_scheduled_task_creation.yml"
  },
  "tags": [
    "attackexecution",
    "attackpersistence",
    "attackprivilege-escalation",
    "attackt1053005"
  ],
  "author": "Nasreddine Bencherchali (Nextron Systems)",
  "references": [
    "https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4698"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2022-12-05",
  "modified": "2022-12-07",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_eid: {\"EventID\": 4698}\nselection_paths: {\"TaskContent|contains\": [\"\\\\AppData\\\\Local\\\\Temp\\\\\", \"\\\\AppData\\\\Roaming\\\\\", \"\\\\Users\\\\Public\\\\\", \"\\\\WINDOWS\\\\Temp\\\\\", \"C:\\\\Temp\\\\\", \"\\\\Desktop\\\\\", \"\\\\Downloads\\\\\", \"\\\\Temporary Internet\", \"C:\\\\ProgramData\\\\\", \"C:\\\\Perflogs\\\\\"]}\nselection_commands: {\"TaskContent|contains\": [\"regsvr32\", \"rundll32\", \"cmd.exe</Command>\", \"cmd</Command>\", \"<Arguments>/c \", \"<Arguments>/k \", \"<Arguments>/r \", \"powershell\", \"pwsh\", \"mshta\", \"wscript\", \"cscript\", \"certutil\", \"bitsadmin\", \"bash.exe\", \"bash \", \"scrcons\", \"wmic \", \"wmic.exe\", \"forfiles\", \"scriptrunner\", \"hh.exe\"]}",
    "condition": "all of selection_*",
    "fields": []
  },
  "falsepositives": [
    "Unknown"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "security"
  ]
}