# rules/sigma/46490193-1b22-4c29-bdd6-5bf63907216f.json
{
  "id": "46490193-1b22-4c29-bdd6-5bf63907216f",
  "name": "VBScript Payload Stored in Registry",
  "description": "Detects VBScript content stored into registry keys as seen being used by UNC2452 group",
  "source": {
    "type": "sigma",
    "id": "46490193-1b22-4c29-bdd6-5bf63907216f",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/registry_set_vbs_payload_stored.yml",
    "file_path": "tmp/repos/sigma/rules/windows/registry/registry_set/registry_set_vbs_payload_stored.yml"
  },
  "tags": [
    "attackpersistence",
    "attackt1547001"
  ],
  "author": "Florian Roth (Nextron Systems)",
  "references": [
    "https://www.microsoft.com/security/blog/2021/03/04/goldmax-goldfinder-sibot-analyzing-nobelium-malware/"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2021-03-05",
  "modified": "2023-08-17",
  "mitre": {
    "tactics": [
      "TA1547"
    ],
    "techniques": [
      "T001"
    ]
  },
  "detection": {
    "query": "selection: {\"TargetObject|contains\": \"Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\", \"Details|contains\": [\"vbscript:\", \"jscript:\", \"mshtml,\", \"RunHTMLApplication\", \"Execute(\", \"CreateObject\", \"window.close\"]}\nfilter: {\"TargetObject|contains\": \"Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\"}\nfilter_dotnet: {\"Image|endswith\": \"\\\\msiexec.exe\", \"TargetObject|contains\": \"\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Installer\\\\UserData\\\\\", \"Details|contains\": [\"\\\\Microsoft.NET\\\\Primary Interop Assemblies\\\\Microsoft.mshtml.dll\", \"<\\\\Microsoft.mshtml,fileVersion=\", \"_mshtml_dll_\", \"<\\\\Microsoft.mshtml,culture=\"]}",
    "condition": "selection and not 1 of filter*",
    "fields": []
  },
  "falsepositives": [
    "Unknown"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "registry_set"
  ]
}