# rules/sigma/4720b7df-40c3-48fd-bbdf-fd4b3c464f0d.json
{
  "id": "4720b7df-40c3-48fd-bbdf-fd4b3c464f0d",
  "name": "Scheduled TaskCache Change by Uncommon Program",
  "description": "Monitor the creation of a new key under 'TaskCache' when a new scheduled task is registered by a process that is not svchost.exe, which is suspicious",
  "source": {
    "type": "sigma",
    "id": "4720b7df-40c3-48fd-bbdf-fd4b3c464f0d",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/registry_set_taskcache_entry.yml",
    "file_path": "tmp/repos/sigma/rules/windows/registry/registry_set/registry_set_taskcache_entry.yml"
  },
  "tags": [
    "attackpersistence",
    "attackt1053",
    "attackt1053005"
  ],
  "author": "Syed Hasan (@syedhasan009)",
  "references": [
    "https://thedfirreport.com/2021/03/29/sodinokibi-aka-revil-ransomware/",
    "https://labs.f-secure.com/blog/scheduled-task-tampering/"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2021-06-18",
  "modified": "2023-08-17",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection: {\"TargetObject|contains\": \"SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\\"}\nfilter: {\"TargetObject|contains\": [\"Microsoft\\\\Windows\\\\UpdateOrchestrator\", \"Microsoft\\\\Windows\\\\SoftwareProtectionPlatform\\\\SvcRestartTask\\\\Index\", \"Microsoft\\\\Windows\\\\Flighting\\\\OneSettings\\\\RefreshCache\\\\Index\"]}\nfilter_tiworker: {\"Image|startswith\": \"C:\\\\Windows\\\\\", \"Image|endswith\": \"\\\\TiWorker.exe\"}\nfilter_svchost: {\"Image\": \"C:\\\\WINDOWS\\\\system32\\\\svchost.exe\"}\nfilter_ngen: {\"Image|startswith\": \"C:\\\\Windows\\\\Microsoft.NET\\\\Framework\", \"Image|endswith\": \"\\\\ngen.exe\", \"TargetObject|contains\": [\"\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tasks\\\\{B66B135D-DA06-4FC4-95F8-7458E1D10129}\", \"\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tree\\\\Microsoft\\\\Windows\\\\.NET Framework\\\\.NET Framework NGEN\"]}\nfilter_office_click_to_run: {\"Image\": [\"C:\\\\Program Files\\\\Microsoft Office\\\\root\\\\Integration\\\\Integrator.exe\", \"C:\\\\Program Files (x86)\\\\Microsoft Office\\\\root\\\\Integration\\\\Integrator.exe\"]}\nfilter_msiexec: {\"Image\": \"C:\\\\Windows\\\\System32\\\\msiexec.exe\"}\nfilter_dropbox_updater: {\"Image\": [\"C:\\\\Program Files (x86)\\\\Dropbox\\\\Update\\\\DropboxUpdate.exe\", \"C:\\\\Program Files\\\\Dropbox\\\\Update\\\\DropboxUpdate.exe\"]}\nfilter_explorer: {\"Image\": \"C:\\\\Windows\\\\explorer.exe\", \"TargetObject|contains\": \"\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tree\\\\Microsoft\\\\Windows\\\\PLA\\\\Server Manager Performance Monitor\\\\\"}\nfilter_system: {\"Image\": \"System\"}",
    "condition": "selection and not 1 of filter*",
    "fields": []
  },
  "falsepositives": [
    "Unknown"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "registry_set"
  ]
}