# rules/sigma/61a7697c-cb79-42a8-a2ff-5f0cdfae0130.json
{
  "id": "61a7697c-cb79-42a8-a2ff-5f0cdfae0130",
  "name": "Potential CobaltStrike Service Installations - Registry",
  "description": "Detects known malicious service installs that appear in cases in which a Cobalt Strike beacon elevates privileges or lateral movement.\n",
  "source": {
    "type": "sigma",
    "id": "61a7697c-cb79-42a8-a2ff-5f0cdfae0130",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/registry_set_cobaltstrike_service_installs.yml",
    "file_path": "tmp/repos/sigma/rules/windows/registry/registry_set/registry_set_cobaltstrike_service_installs.yml"
  },
  "tags": [
    "attackexecution",
    "attacklateral-movement",
    "attackprivilege-escalation",
    "attackt1021002",
    "attackt1543003",
    "attackt1569002"
  ],
  "author": "Wojciech Lesicki",
  "references": [
    "https://www.sans.org/webcasts/tech-tuesday-workshop-cobalt-strike-detection-log-analysis-119395"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2021-06-29",
  "modified": "2024-03-25",
  "mitre": {
    "tactics": [
      "TA1021",
      "TA1543",
      "TA1569"
    ],
    "techniques": [
      "T002",
      "T003"
    ]
  },
  "detection": {
    "query": "selection_key: [{\"TargetObject|contains\": \"\\\\System\\\\CurrentControlSet\\\\Services\"}, {\"TargetObject|contains|all\": [\"\\\\System\\\\ControlSet\", \"\\\\Services\"]}]\nselection_details: [{\"Details|contains|all\": [\"ADMIN$\", \".exe\"]}, {\"Details|contains|all\": [\"%COMSPEC%\", \"start\", \"powershell\"]}]",
    "condition": "all of selection_*",
    "fields": []
  },
  "falsepositives": [
    "Unlikely"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "registry_set"
  ]
}