# rules/sigma/6b14bac8-3e3a-4324-8109-42f0546a347f.json
{
  "id": "6b14bac8-3e3a-4324-8109-42f0546a347f",
  "name": "Scheduled Cron Task/Job - Linux",
  "description": "Detects abuse of the cron utility to perform task scheduling for initial or recurring execution of malicious code. Detection will focus on crontab jobs uploaded from the tmp folder.",
  "source": {
    "type": "sigma",
    "id": "6b14bac8-3e3a-4324-8109-42f0546a347f",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_lnx_schedule_task_job_cron.yml",
    "file_path": "tmp/repos/sigma/rules/linux/process_creation/proc_creation_lnx_schedule_task_job_cron.yml"
  },
  "tags": [
    "attackexecution",
    "attackpersistence",
    "attackprivilege-escalation",
    "attackt1053003"
  ],
  "author": "Alejandro Ortuno, oscd.community",
  "references": [
    "https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md"
  ],
  "severity": "medium",
  "type": "sigma",
  "status": "test",
  "created": "2020-10-06",
  "modified": "2022-11-27",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T003"
    ]
  },
  "detection": {
    "query": "selection: {\"Image|endswith\": \"crontab\", \"CommandLine|contains\": \"/tmp/\"}",
    "condition": "selection",
    "fields": []
  },
  "falsepositives": [
    "Legitimate administration activities"
  ],
  "level": "medium",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "linux",
    "process_creation"
  ]
}