# rules/sigma/6c4e2f43-d94d-4ead-b64d-97e53fa2bd05.json
{
  "id": "6c4e2f43-d94d-4ead-b64d-97e53fa2bd05",
  "name": "Persistence Via Cron Files",
  "description": "Detects creation of cron file or files in Cron directories which could indicates potential persistence.",
  "source": {
    "type": "sigma",
    "id": "6c4e2f43-d94d-4ead-b64d-97e53fa2bd05",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/file_event_lnx_persistence_cron_files.yml",
    "file_path": "tmp/repos/sigma/rules/linux/file_event/file_event_lnx_persistence_cron_files.yml"
  },
  "tags": [
    "attackpersistence",
    "attackt1053003"
  ],
  "author": "Roberto Rodriguez (Cyb3rWard0g), OTR (Open Threat Research), MSTIC",
  "references": [
    "https://github.com/microsoft/MSTIC-Sysmon/blob/f1477c0512b0747c1455283069c21faec758e29d/linux/configs/attack-based/persistence/T1053.003_Cron_Activity.xml"
  ],
  "severity": "medium",
  "type": "sigma",
  "status": "test",
  "created": "2021-10-15",
  "modified": "2022-12-31",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T003"
    ]
  },
  "detection": {
    "query": "selection1: {\"TargetFilename|startswith\": [\"/etc/cron.d/\", \"/etc/cron.daily/\", \"/etc/cron.hourly/\", \"/etc/cron.monthly/\", \"/etc/cron.weekly/\", \"/var/spool/cron/crontabs/\"]}\nselection2: {\"TargetFilename|contains\": [\"/etc/cron.allow\", \"/etc/cron.deny\", \"/etc/crontab\"]}",
    "condition": "1 of selection*",
    "fields": []
  },
  "falsepositives": [
    "Any legitimate cron file."
  ],
  "level": "medium",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "linux",
    "file_event"
  ]
}