# rules/sigma/7a02e22e-b885-4404-b38b-1ddc7e65258a.json
{
  "id": "7a02e22e-b885-4404-b38b-1ddc7e65258a",
  "name": "Suspicious Schtasks Schedule Type With High Privileges",
  "description": "Detects scheduled task creations or modification to be run with high privileges on a suspicious schedule type",
  "source": {
    "type": "sigma",
    "id": "7a02e22e-b885-4404-b38b-1ddc7e65258a",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_schedule_type_system.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_schedule_type_system.yml"
  },
  "tags": [
    "attackexecution",
    "attackt1053005"
  ],
  "author": "Nasreddine Bencherchali (Nextron Systems)",
  "references": [
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change",
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create"
  ],
  "severity": "medium",
  "type": "sigma",
  "status": "test",
  "created": "2022-08-31",
  "modified": "2022-08-31",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_img: [{\"Image|endswith\": \"\\\\schtasks.exe\"}, {\"OriginalFileName\": \"schtasks.exe\"}]\nselection_time: {\"CommandLine|contains\": [\" ONLOGON \", \" ONSTART \", \" ONCE \", \" ONIDLE \"]}\nselection_privs: {\"CommandLine|contains\": [\"NT AUT\", \" SYSTEM\", \"HIGHEST\"]}",
    "condition": "all of selection_*",
    "fields": []
  },
  "falsepositives": [
    "Some installers were seen using this method of creation unfortunately. Filter them in your environment"
  ],
  "level": "medium",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}