# rules/sigma/92626ddd-662c-49e3-ac59-f6535f12d189.json
{
  "id": "92626ddd-662c-49e3-ac59-f6535f12d189",
  "name": "Scheduled Task Creation Via Schtasks.EXE",
  "description": "Detects the creation of scheduled tasks by user accounts via the \"schtasks\" utility.",
  "source": {
    "type": "sigma",
    "id": "92626ddd-662c-49e3-ac59-f6535f12d189",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_creation.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_creation.yml"
  },
  "tags": [
    "attackexecution",
    "attackpersistence",
    "attackprivilege-escalation",
    "attacks0111",
    "attackt1053005",
    "car2013-08-001",
    "stp1u"
  ],
  "author": "Florian Roth (Nextron Systems)",
  "references": [
    "https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create"
  ],
  "severity": "low",
  "type": "sigma",
  "status": "test",
  "created": "2019-01-16",
  "modified": "2024-01-18",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection: {\"Image|endswith\": \"\\\\schtasks.exe\", \"CommandLine|contains\": \" /create \"}\nfilter_main_system_user: {\"User|contains\": [\"AUTHORI\", \"AUTORI\"]}",
    "condition": "selection and not 1 of filter_main_*",
    "fields": []
  },
  "falsepositives": [
    "Administrative activity",
    "Software installation"
  ],
  "level": "low",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}