# rules/sigma/a8f29a7b-b137-4446-80a0-b804272f3da2.json
{
  "id": "a8f29a7b-b137-4446-80a0-b804272f3da2",
  "name": "Persistence and Execution at Scale via GPO Scheduled Task",
  "description": "Detect lateral movement using GPO scheduled task, usually used to deploy ransomware at scale",
  "source": {
    "type": "sigma",
    "id": "a8f29a7b-b137-4446-80a0-b804272f3da2",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/win_security_gpo_scheduledtasks.yml",
    "file_path": "tmp/repos/sigma/rules/windows/builtin/security/win_security_gpo_scheduledtasks.yml"
  },
  "tags": [
    "attacklateral-movement",
    "attackpersistence",
    "attackt1053005"
  ],
  "author": "Samir Bousseaden",
  "references": [
    "https://twitter.com/menasec1/status/1106899890377052160",
    "https://www.secureworks.com/blog/ransomware-as-a-distraction",
    "https://www.elastic.co/guide/en/security/7.17/prebuilt-rule-0-16-1-scheduled-task-execution-at-scale-via-gpo.html"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2019-04-03",
  "modified": "2024-09-04",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_5136: {\"EventID\": 5136, \"AttributeLDAPDisplayName\": [\"gPCMachineExtensionNames\", \"gPCUserExtensionNames\"], \"AttributeValue|contains\": [\"CAB54552-DEEA-4691-817E-ED4A4D1AFC72\", \"AADCED64-746C-4633-A97C-D61349046527\"]}\nselection_5145: {\"EventID\": 5145, \"ShareName|endswith\": \"\\\\SYSVOL\", \"RelativeTargetName|endswith\": \"ScheduledTasks.xml\", \"AccessList|contains\": [\"WriteData\", \"%%4417\"]}",
    "condition": "1 of selection_*",
    "fields": []
  },
  "falsepositives": [
    "If the source IP is not localhost then it's super suspicious, better to monitor both local and remote changes to GPO scheduled tasks."
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "security"
  ]
}