# rules/sigma/f2c64357-b1d2-41b7-849f-34d2682c0fad.json
{
  "id": "f2c64357-b1d2-41b7-849f-34d2682c0fad",
  "name": "Suspicious Command Patterns In Scheduled Task Creation",
  "description": "Detects scheduled task creation using \"schtasks\" that contain potentially suspicious or uncommon commands",
  "source": {
    "type": "sigma",
    "id": "f2c64357-b1d2-41b7-849f-34d2682c0fad",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_susp_pattern.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_susp_pattern.yml"
  },
  "tags": [
    "attackexecution",
    "attackt1053005"
  ],
  "author": "Florian Roth (Nextron Systems)",
  "references": [
    "https://app.any.run/tasks/512c1352-6380-4436-b27d-bb62f0c020d6/",
    "https://twitter.com/RedDrip7/status/1506480588827467785",
    "https://www.ncsc.gov.uk/static-assets/documents/malware-analysis-reports/devil-bait/NCSC-MAR-Devil-Bait.pdf"
  ],
  "severity": "high",
  "type": "sigma",
  "status": "test",
  "created": "2022-02-23",
  "modified": "2024-03-19",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_schtasks: {\"Image|endswith\": \"\\\\schtasks.exe\", \"CommandLine|contains\": \"/Create \"}\nselection_pattern_1: {\"CommandLine|contains\": [\"/sc minute \", \"/ru system \"]}\nselection_pattern_2: {\"CommandLine|contains\": [\"cmd /c\", \"cmd /k\", \"cmd /r\", \"cmd.exe /c \", \"cmd.exe /k \", \"cmd.exe /r \"]}\nselection_uncommon: {\"CommandLine|contains\": [\" -decode \", \" -enc \", \" -w hidden \", \" bypass \", \" IEX\", \".DownloadData\", \".DownloadFile\", \".DownloadString\", \"/c start /min \", \"FromBase64String\", \"mshta http\", \"mshta.exe http\"]}\nselection_anomaly_1: {\"CommandLine|contains\": [\":\\\\ProgramData\\\\\", \":\\\\Temp\\\\\", \":\\\\Tmp\\\\\", \":\\\\Users\\\\Public\\\\\", \":\\\\Windows\\\\Temp\\\\\", \"\\\\AppData\\\\\", \"%AppData%\", \"%Temp%\", \"%tmp%\"]}\nselection_anomaly_2: {\"CommandLine|contains\": [\"cscript\", \"curl\", \"wscript\"]}",
    "condition": "selection_schtasks and ( all of selection_pattern_* or selection_uncommon or all of selection_anomaly_* )",
    "fields": []
  },
  "falsepositives": [
    "Software installers that run from temporary folders and also install scheduled tasks are expected to generate some false positives"
  ],
  "level": "high",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}