# rules/sigma/f548a603-c9f2-4c89-b511-b089f7e94549.json
{
  "id": "f548a603-c9f2-4c89-b511-b089f7e94549",
  "name": "Potential Persistence Via Microsoft Compatibility Appraiser",
  "description": "Detects manual execution of the \"Microsoft Compatibility Appraiser\" task via schtasks.\nIn order to trigger persistence stored in the \"\\AppCompatFlags\\TelemetryController\" registry key.\n",
  "source": {
    "type": "sigma",
    "id": "f548a603-c9f2-4c89-b511-b089f7e94549",
    "url": "https://github.com/SigmaHQ/sigma/blob/master/proc_creation_win_schtasks_persistence_windows_telemetry.yml",
    "file_path": "tmp/repos/sigma/rules/windows/process_creation/proc_creation_win_schtasks_persistence_windows_telemetry.yml"
  },
  "tags": [
    "attackpersistence",
    "attackt1053005"
  ],
  "author": "Sreeman",
  "references": [
    "https://www.trustedsec.com/blog/abusing-windows-telemetry-for-persistence/"
  ],
  "severity": "medium",
  "type": "sigma",
  "status": "test",
  "created": "2020-09-29",
  "modified": "2023-02-10",
  "mitre": {
    "tactics": [
      "TA1053"
    ],
    "techniques": [
      "T005"
    ]
  },
  "detection": {
    "query": "selection_img: [{\"Image|endswith\": \"\\\\schtasks.exe\"}, {\"OriginalFileName\": \"schtasks.exe\"}]\nselection_cli: {\"CommandLine|contains|all\": [\"run \", \"\\\\Application Experience\\\\Microsoft Compatibility Appraiser\"]}",
    "condition": "all of selection_*",
    "fields": []
  },
  "falsepositives": [
    "Unknown"
  ],
  "level": "medium",
  "rule_format": "standard",
  "platforms": [],
  "data_sources": [
    "windows",
    "process_creation"
  ]
}